# ansible/systemd/bittrex.service
[Unit]
Description=Bittrex trader

[Service]
ExecStart={{project_dir}}/bin/bittrex_trader --config {{config}}
