# ansible/systemd/trader.service
[Unit]
Description=Bittrex trader

[Service]
ExecStart={{project_dir}}/bin/trader --config {{config}}
