# ansible/systemd/trader.service
[Unit]
Description=Bittrex trader

[Service]
WorkingDirectory={{project_dir}}
ExecStart={{project_dir}}/bin/trader --config {{config}}
